# polymarket-bot.service
[Unit]
Description=Polymarket HFT Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket
EnvironmentFile=/home/ubuntu/polymarket/.env
ExecStart=/home/ubuntu/polymarket/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
